# rootfs/etc/systemd/system/gps-logger.service
[Unit]
Description=Log nmea from gpsd
After=gpsd.service
Wants=gpsd.service

[Service]
ExecStart=/usr/local/bin/gps-logger.sh

[Install]
WantedBy=multi-user.target
